# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn --timeout=3600 app:app -b localhost:5000
directory=/deploy/app
